# systemd/user/wallpaper-cycle.service
[Unit]
Description=Cycle wallpapers with Hyprpaper

[Service]
ExecStart=%h/bin/cycle-wallpapers.sh
